# misc-utils/uuidd.socket.in
[Unit]
Description=UUID daemon activation socket

[Socket]
ListenStream=@runstatedir@/uuidd/request

[Install]
WantedBy=sockets.target
