# sys-utils/fstrim.service.in
[Unit]
Description=Discard unused blocks on filesystems from /etc/fstab

[Service]
Type=oneshot
ExecStart=@sbindir@/fstrim -Av
