# sys-utils/fstrim.service.in
[Unit]
Description=Discard unused blocks

[Service]
Type=oneshot
ExecStart=@sbindir@/fstrim -a
