# sys-utils/fstrim.service.in
[Unit]
Description=Discard unused blocks on filesystems from /etc/fstab
Documentation=man:fstrim(8)

[Service]
Type=oneshot
ExecStart=@sbindir@/fstrim -Av
