# example.files/inittab
# inittab
# Format:
# ttyline:termcap-entry:getty-command
tty1:linux:/sbin/getty 9600 tty1
tty2:linux:/sbin/getty 9600 tty2
tty3:linux:/sbin/getty 9600 tty3
tty4:linux:/sbin/getty 9600 tty4
tty5:linux:/sbin/getty 9600 tty5
tty6:linux:/sbin/getty 9600 tty6
ttyS1:vt100:/sbin/getty -Lh 19200 ttyS1
